# infra/supervisor.conf
; MFTC Supervisor Config
; Place at /etc/supervisor/conf.d/mftc.conf

[program:mftc-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/mftc/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600
user=www-data
numprocs=3
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/mftc-queue.log
stderr_logfile=/var/log/mftc-queue-error.log
stopwaitsecs=3600

; Scheduler menjalankan job MFTC + spatie/laravel-backup (backup:clean 01:30, backup:run 02:30).
; Tidak perlu program supervisor terpisah untuk backup — cukup pastikan mftc-schedule aktif.
[program:mftc-schedule]
command=bash -c 'while true; do php /var/www/mftc/artisan schedule:run; sleep 60; done'
user=www-data
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/var/log/mftc-schedule.log
stderr_logfile=/var/log/mftc-schedule-error.log
